# gssiwmd-default-settings/files/etc/uci-defaults/zzzz-gssiwmd-settings
#!/bin/sh
uci set luci.main.lang=zh_cn
uci commit luci

uci set system.@system[0].timezone=CST-8
uci set system.@system[0].zonename=Asia/Shanghai
uci commit system

uci set fstab.@global[0].anon_mount=1
uci commit fstab
DISTRIB_ID='ImmortalWrt'

sed -i '/DISTRIB_ID/d' /etc/openwrt_release
echo "DISTRIB_ID='HomeRouter'" >> /etc/openwrt_release

sed -i '/DISTRIB_REVISION/d' /etc/openwrt_release
echo "DISTRIB_REVISION='(Compiled by WMD)'" >> /etc/openwrt_release
sed -i '/DISTRIB_DESCRIPTION/d' /etc/openwrt_release
echo "DISTRIB_DESCRIPTION='OpenWrt '" >> /etc/openwrt_release
sed -i '/DISTRIB_LASTMODIFIELD/d' /etc/openwrt_release
echo "DISTRIB_LASTMODIFIELD='$(date "+%Y-%m-%d %H:%M:%S")'" >> /etc/openwrt_release
# Replace Hostname
sed -i "s/hostname .*/hostname 'HomeRouter'/" /etc/config/system
# Default password of root
sed -i 's/root:.*/root:$1$Ds8BJ34B$Lm19IlaiTFvnFFLYDvEU6.:19172:0:99999:7:::/' /etc/shadow
exit 0
