# gssiwmd-default-settings/files/etc/uci-defaults/zzzz-gssiwmd-settings
#!/bin/sh
uci set luci.main.lang=zh_cn
uci commit luci

uci set system.@system[0].timezone=CST-8
uci set system.@system[0].zonename=Asia/Shanghai
uci commit system

uci set fstab.@global[0].anon_mount=1
uci commit fstab

sed -i '/DISTRIB_ID/d' /etc/openwrt_release
echo "DISTRIB_ID='gssiRouter '" >> /etc/openwrt_release
sed -i '/DISTRIB_REVISION/d' /etc/openwrt_release
echo "DISTRIB_REVISION='$(date "+%Y-%m-%d")'" >> /etc/openwrt_release
sed -i '/DISTRIB_DESCRIPTION/d' /etc/openwrt_release
echo "DISTRIB_DESCRIPTION='gssiRouter '" >> /etc/openwrt_release
sed -i '/DISTRIB_LASTMODIFIELD/d' /etc/openwrt_release
echo "DISTRIB_LASTMODIFIELD='Compiled by Wangmingde $(date "+%Y-%m-%d %H:%M:%S")'" >> /etc/openwrt_release
# Replace system
# sed -i "s/hostname .*/hostname 'gssiRouter'/" /etc/config/system

#Default password of root
sed -i 's/root:.*/root:$1$Ds8BJ34B$Lm19IlaiTFvnFFLYDvEU6.:19172:0:99999:7:::/' /etc/shadow

mv /etc/gssiwmd_banner /etc/banner

mv /etc/gssiwmd_system /etc/config/system

mv /etc/gssiwmd_sysinfo.sh /etc/profile.d/sysinfo.sh



exit 0
